# src/console/glimpse-probe.conf
description	"Runs the glimpse probe"

start on (local-filesystems and net-device-up IFACE!=lo)
stop on runlevel [016]

#respawn

script
        sudo -H -u glimpse glimpse-console
end script
