# src/console/glimpse-probe.conf
description	"Runs the glimpse probe"

start on runlevel [2345]
stop on runlevel [016]

#respawn

script
        sudo -H -u glimpse /opt/glimpse/glimpse-console
end script
