# src/console/glimpse.service
[Unit]
Description=Glimpse measurement client
After=network.target

[Service]
User=glimpse
ExecStart=/opt/glimpse/glimpse-console

[Install]
WantedBy=multi-user.target
